# dist/fidorium.initd
#!/sbin/openrc-run

description="FIDO2 hardware authenticator daemon"
command="/usr/bin/fidorium"
command_args="${FIDORIUM_ARGS:--v}"
command_background=yes
pidfile="/run/fidorium.pid"
output_log="/var/log/fidorium.log"
error_log="/var/log/fidorium.log"

depend() {
    after udev
}
