# dist/fidorium.initd
#!/sbin/openrc-run

description="FIDO2 hardware authenticator daemon"
command="/usr/bin/fidorium"
command_args="${FIDORIUM_ARGS:--v}"
command_background=yes
pidfile="${XDG_RUNTIME_DIR}/fidorium.pid"
output_log="${XDG_RUNTIME_DIR}/fidorium.log"
error_log="${XDG_RUNTIME_DIR}/fidorium.log"

depend() {
    after udev
}
